# containers/build/entrypoint.service
[Unit]
Description=Checks the server status

[Service]
User=root
ExecStart=/bin/bash /root/entrypoint.sh

[Install]
WantedBy=multi-user.target